# etc/profile.d/alias-dotfiles-man
#!/bin/bash

alias dotm="~/usr/bin/dotfiles-man"
alias dman="~/usr/bin/dotfiles-man"
alias dotmr="~/usr/bin/dotfiles-man -r"
alias dmanr="~/usr/bin/dotfiles-man -r"
