# etc/profile.d/alias-dotfiles-man
#!/bin/bash

alias dotm="~/usr/bin/dotfiles-man"
