# etc/profile.d/dotfiles-man
#!/bin/sh

alias dotm="~/usr/bin/dotfiles-man"
