# etc/profile.d/dotfiles-man
#!/bin/sh

alias dm="~/usr/bin/dotfiles-man"
