# etc/profile.d/dotfiles-man
#!/bin/sh

alias dfm="~/usr/bin/dotfiles-man"
